# docker/s6-rc.d/api/run
#!/command/with-contenv sh
# The Rust API: reads SETS_DIRECTORY (data/sets) relative to /app and binds the
# loopback port the SvelteKit server proxies to (BIND_ADDRESS).
cd /app
exec murmelbahn-web
