# docker/supervisord.conf
# supervisord runs the Rust API and the SvelteKit server in one container. Unlike
# s6-overlay it does not need to be PID 1, so it works under Fly Machines' init
# (which runs the entrypoint as a child process). Program output goes to the
# container's stdout/stderr so it shows up in `fly logs`.
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:api]
command=murmelbahn-web
directory=/app
autorestart=true
startretries=1000000
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:web]
command=node build
directory=/app/frontend
autorestart=true
startretries=1000000
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
